# hosts/networking-and-dns/templates/etc/systemd/system/networking.service.d/auto-restart.conf

# {{ ansible_managed }}

[Service]
Restart=on-failure
RestartSec=5s

[Unit]
StartLimitBurst=60
StartLimitIntervalSec=30
